# request-tracker/3/etc/service/rt/run
#!/bin/sh
exec 2>&1

set -e

exec /sbin/setuser www-data /opt/rt3/bin/fastcgi_server -s /opt/rt3/var/fastcgi.sock